# Actividad_04/actividad4.service
[Unit]
Description=Actividad4

[Service]
Type=simple
ExecStart=/usr/bin/bash ./saludo.sh

[Install]
WantedBy=multi-user.target
